# vless-gateway.service
[Unit]
Description=VLESS Self-Healing Gateway
After=network.target

[Service]
# Укажите пользователя, от которого работает сервис
User=developer
Group=developer
# Укажите полный путь к директории проекта
WorkingDirectory=/home/developer/vless-checker
# Укажите полный путь к python3 и вашему скрипту
ExecStart=/home/developer/bin/python3 /home/developer/vless-checker/gateway_manager.py
Restart=always

[Install]
WantedBy=multi-user.target